# ledtest.service
[Unit]

Description=Led Service Test

After=network.target

[Service]

Type=simple
#Type=forking
#PIDFile=/var/run/ledtest.pid

WorkingDirectory=/home/pi/apps/systemd/
ExecStart=/usr/bin/python /home/pi/apps/systemd/ledtest.py
KillMode=process
#RemainAfterExit=true

[Install]

WantedBy=multi-user.target
